# docker/supervisord.conf
[program:uwsgi]
command=/usr/local/bin/uwsgi --ini /app/app.ini
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx
startretries=0
